# upbit-bot.service
[Unit]
Description=Upbit Signal Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bot
ExecStartPre=/bin/rm -f /home/ubuntu/bot/.study_running /home/ubuntu/bot/.study_lock /home/ubuntu/bot/.study_done
ExecStart=/usr/bin/python3 /home/ubuntu/bot/bot.py
EnvironmentFile=/home/ubuntu/bot/.env
Restart=always
RestartSec=60
StandardOutput=append:/tmp/bot_stdout.log
StandardError=append:/tmp/bot_stdout.log

[Install]
WantedBy=multi-user.target
