# upbit-bot.service
[Unit]
Description=Upbit Signal Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bots
ExecStartPre=/bin/rm -f /home/ubuntu/bots/.study_running /home/ubuntu/bots/.study_lock /home/ubuntu/bots/.study_done
ExecStart=/usr/bin/python3 /home/ubuntu/bots/bot.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/tmp/bot_stdout.log
StandardError=append:/tmp/bot_stdout.log

[Install]
WantedBy=multi-user.target
